# apps/backend/systemd/Openoverlaybackend.service
[Unit]
Description=OpenOverlay backend API and WebSocket server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=skylarenns
Group=skylarenns
WorkingDirectory=/home/skylarenns/Documents/GitHub/OpenOverlay/apps/backend
EnvironmentFile=/etc/openoverlaybackend.env
Environment=SELF_UPDATE_ENABLED=true
Environment=SELF_UPDATE_REPO_DIR=/home/skylarenns/Documents/GitHub/OpenOverlay
Environment=SELF_UPDATE_REMOTE=origin
Environment=SELF_UPDATE_BRANCH=main
Environment=SELF_UPDATE_INTERVAL_MS=60000
Environment=GATEWAY_BACKEND_PORTS=8735,8736
Environment=GATEWAY_RELEASE_DIR=/home/skylarenns/Documents/GitHub/OpenOverlay/releases
Environment=OPENOVERLAY_GATEWAY_ENTRYPOINT=1
ExecStart=/usr/bin/node dist/gateway.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
